# dist/dracut/modules.d/99attestation-agent/attestation-agent-platform-detect.service
[Unit]
Description=Attestation-Agent Detect Platform Attributes
DefaultDependencies=no
ConditionPathExists=/etc/initrd-release
After=sysinit.target

[Service]
Type=oneshot
ExecStart=bash /usr/bin/attestation-agent-platform-detect.sh
RemainAfterExit=true
