# dist/dracut/modules.d/99attestation-agent/attestation-agent.service
[Unit]
Description=Attestation-Agent
Documentation=https://confidentialcontainers.org
DefaultDependencies=no
ConditionPathExists=/etc/initrd-release
After=sysinit.target
Before=basic.target

[Service]
ExecStartPre=-/usr/sbin/modprobe tdx-guest
ExecStart=/usr/local/bin/attestation-agent --config-file /etc/trustiflux/attestation-agent.toml -a unix:///run/confidential-containers/attestation-agent/attestation-agent.sock
Environment=RUST_LOG=debug
StandardOutput=journal+console
StandardError=journal+console
Restart=always
RestartSec=5
Delegate=yes
KillMode=process
OOMScoreAdjust=-999
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity

[Install]
WantedBy=basic.target