# dist/rpm/attestation-agent.service
[Unit]
Description=Attestation-Agent
Documentation=https://confidentialcontainers.org
After=network.target

[Service]
ExecStart=/usr/bin/attestation-agent --config-file /etc/trustiflux/attestation-agent.toml -a unix:///run/confidential-containers/attestation-agent/attestation-agent.sock
Environment=RUST_LOG=debug
Restart=always
RestartSec=5
Delegate=yes
KillMode=process
OOMScoreAdjust=-999
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity

[Install]
WantedBy=multi-user.target