# dist/rpm/trustiflux-api-server.service
[Unit]
Description=Trustiflux API Server (REST bridge for attestation-agent and confidential-data-hub)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/trustiflux-api-server --config /etc/trustiflux/trustiflux-api-server.toml
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target
