# dist/rpm/trustiflux-api-server.service
[Unit]
Description=Trustiflux API Server (REST bridge for attestation-agent and confidential-data-hub)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/trustiflux-api-server --bind 0.0.0.0:8006 --features all --aa_addr unix:///run/confidential-containers/attestation-agent/attestation-agent.sock --cdh_addr unix:///run/confidential-containers/cdh.sock
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target

